# infra/iaas/webhook.service
[Unit]
Description=Webhook pull all webs
After=network.target

[Service]
ExecStart=/pm/infra/iaas/pull.sh

[Install]
WantedBy=multi-user.target
